# pacman/pacman-sync-db.service
[Unit]
Description=Synchronize Pacman Database

[Service]
ExecStart=/usr/bin/pacman -Sy
